# crontab
0 * * * * source /Users/tyleranderson/PyCharmProjects/python_quickstart/jobs/backup.sh

